# main/crontab
0 0,12 * * * python3 /code/main/orders_to_excel.py >> /var/log/cron.log 2>&1
